# tarsnap/tarsnap.timer
[Unit]
Description=Tarsnap backup

[Timer]
OnCalendar=monday 12:00
Persistent=true

[Install]
WantedBy=timers.target
