# tarsnap/tarsnap.timer
[Unit]
Description=Tarsnap backup

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
